# supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod

[program:apache2]
command=/usr/sbin/apache2 -DFOREGROUND
